# packages/postgresml-dashboard/DEBIAN/prerm
#!/bin/bash

systemctl stop pgml-dashboard
systemctl disable pgml-dashboard
